# api/gitlab/cronjob/crontab
* * * * *  ./api/gitlab/cronjob/scripts/monitoring_domain.py >> ~/cron.log 2>&1
