# api/gitlab/cronjob/crontab
* * * * * ./ada/cronjob/scripts/monitoring_domain.py
